# gnos_transfer_monitor/check_outages.cron
#!/bin/bash

#source /home/ubuntu/perl5/perlbrew/etc/bashrc

cd /home/ubuntu/gitroot/pancancer-sandbox/gnos_transfer_monitor/

#export PATH=$PATH:/usr/local/bin

./check_outages.pl >> outage_check.log 2>&1
