# gnos_transfer_monitor/generate_html_report.cron
#!/bin/bash

source /home/ubuntu/perl5/perlbrew/etc/bashrc

cd /home/ubuntu/gitroot/pancancer-sandbox/gnos_transfer_monitor/

export PATH=$PATH:/usr/local/bin

perl generate_html_report.pl --template map.tt > /var/www/dev/transfers.html
perl generate_html_report.pl --template map.tt > /var/www/transfers.html
